# cookiecutter-templates/cookiecutter-dioptra-deployment/{{cookiecutter.__project_slug}}/systemd/dioptra-backup.timer
[Unit]
Description=Dioptra Volume Backup timer
RefuseManualStart=no
RefuseManualStop=no

[Timer]
Persistent=true
OnCalendar=*-*-* 06:00:00
Unit=dioptra-backup.service

[Install]
WantedBy=timers.target
